# systemd/telegram-lark-backfill.timer
[Unit]
Description=Run Telegram to Lark backfill every 5 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=30s
Persistent=true

[Install]
WantedBy=timers.target
