# systemd/telegram-lark-ai-review.timer
[Unit]
Description=Run AI Telegram discussion review every 30 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min
AccuracySec=1min
Persistent=true

[Install]
WantedBy=timers.target
