# systemd/telegram-lark-intake.timer
[Unit]
Description=Run Telegram discussion to Lark Base intake every 5 minutes

[Timer]
OnBootSec=3min
OnUnitActiveSec=5min
AccuracySec=30s
Persistent=true

[Install]
WantedBy=timers.target
